# chassis/launchd/com.behalfbot.discord-restart.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.behalfbot.${BOT_NAME}-discord-restart.plist - chassis-rendered template.

  Generated from chassis/launchd/com.behalfbot.discord-restart.plist.template
  by chassis/scripts/bootstrap-customer-scripts.sh.

  DOMAIN: LaunchAgent (user-level, ~/Library/LaunchAgents/), bootstrapped into
  gui/<uid>. This MUST NOT be a LaunchDaemon. See docs/launchd-domains.md.

  WHY (correcting chassis#14):
  #14 promoted this job to a LaunchDaemon on the claim that it "only docker
  execs the chassis container - no GUI/Aqua session needed". That claim was
  false. The job runs ${CUSTOMER_HOME}/scripts/restart-${BOT_NAME}-discord.sh,
  which spawns a HOST tmux session running `claude`. That Claude process, and
  everything it shells out to, needs the user's login keychain.

  A LaunchDaemon runs in launchd's Background session even with UserName set.
  A Background session cannot unlock the user's login keychain without GUI
  interaction, so `security find-generic-password` fails with error 36 ("User
  interaction is not allowed"). Downstream, bw-unlock.sh cannot read the
  Vaultwarden master password and every Vaultwarden-sourced credential dies.

  Two traps that made this hard to see, and that any future change must respect:

    1. tmux runs ONE server per user socket (/tmp/tmux-<uid>/default). Whichever
       job creates that server fixes the launchd session for EVERY session on
       it. A single Background-domain job spawning any tmux session poisons the
       server for all of them. Converting only some of the tmux-spawning jobs
       to agents achieves nothing - it has to be all of them.
    2. The restart/watchdog scripts historically did `tmux kill-session`, not
       `kill-server`. A poisoned server survives every restart forever - the
       reference install stayed poisoned through five weeks of daily 05:00
       restarts. The rendered restart script now detects a non-Aqua-born server
       and kills the SERVER so launchd rebuilds it in the Aqua session. See
       chassis/scripts/templates/restart-discord.sh.template.

  Outage: 2026-06-03 (the #14 daemon sweep) to 2026-07-11 on the reference
  install. Every macOS install shipped since #14 is affected.

  TRADEOFF (do not "fix" this back to a daemon):
  #14's motivation was real - a gui LaunchAgent only loads after GUI login, so
  it does not recover on an unattended reboot of a Mac with auto-login off. On
  macOS you cannot have both login-keychain access and pre-login startup: the
  login keychain is unlocked BY the user login. There is no clever workaround.
  The correct answer is LaunchAgent + auto-login enabled on the install Mac
  (System Settings > Users & Groups > Automatically log in as).

  Template variables (expanded at render time):
    BOT_NAME            logical bot name (e.g. jax, asimov, ozzy)
    CUSTOMER_HOME       per-customer state root (e.g. /Users/<user>/.behalfbot)
    CHASSIS_HOME        chassis git tree root (e.g. /Users/<user>/behalfbot)
    HOMEBREW_PREFIX     /opt/homebrew on Apple Silicon, /usr/local on Intel
    USER                installer's macOS short name
    NODE_BIN_PATH       optional - extra node bin dir, e.g. ${HOMEBREW_PREFIX}/opt/node@22/bin

  No UserName/GroupName keys: a LaunchAgent already runs as the logged-in user.
  Those keys are daemon-only. Paths stay absolute anyway - costs nothing, and
  keeps the plist readable.

  Restart cadence: daily at 05:00 local time. The daily-briefing heartbeat fires
  after the restart, so a fresh tmux session is in place before it runs.

  RunAtLoad=true: also fire once when launchd bootstraps the agent, i.e. at GUI
  login and at activation time. On an install migrating off the old
  LaunchDaemons this is what performs the one-shot tmux-server respawn: the
  restart script sees the Background-born server and rebuilds it under Aqua.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.behalfbot.${BOT_NAME}-discord-restart</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>${CUSTOMER_HOME}/scripts/restart-${BOT_NAME}-discord.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>${CHASSIS_HOME}</string>
    <key>StandardErrorPath</key>
    <string>${CUSTOMER_HOME}/logs/scheduled/discord-restart-stderr.log</string>
    <key>StandardOutPath</key>
    <string>${CUSTOMER_HOME}/logs/scheduled/discord-restart-stdout.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/${USER}/.local/bin:${HOMEBREW_PREFIX}/bin:${NODE_BIN_PATH}:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin</string>
        <key>HOME</key>
        <string>/Users/${USER}</string>
        <key>CUSTOMER_HOME</key>
        <string>${CUSTOMER_HOME}</string>
        <key>CHASSIS_HOME</key>
        <string>${CHASSIS_HOME}</string>
    </dict>
</dict>
</plist>
